# deploy/systemd/dns-update.service
[Unit]
Description=Keep DNS A and AAAA records aligned with current egress IPs
Wants=network-online.target
After=network-online.target
ConditionPathExists=/usr/bin/dns-update
ConditionPathExists=/etc/dns-update/config.json
ConditionPathExists=/etc/dns-update/cloudflare.token

[Service]
Type=oneshot
ExecStart=/usr/bin/dns-update
Environment=DNS_UPDATE_CONFIG=/etc/dns-update/config.json
Environment=DNS_UPDATE_TIMEOUT=2m
Environment=DNS_UPDATE_PROVIDER_CLOUDFLARE_API_TOKEN_FILE=%d/cloudflare.token
EnvironmentFile=-/etc/dns-update/dns-update.env
LoadCredential=cloudflare.token:/etc/dns-update/cloudflare.token
TimeoutStartSec=3m
SyslogIdentifier=dns-update
UMask=0077
DynamicUser=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateDevices=yes
PrivateMounts=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
